# etc/profile.d/entuura-sovereignty.sh
# Entuura Ventures Limited
# Copyright 2016
#
# Author: Steven Uggowitzer
export PATH="/opt/ent-sover/bin:/opt/vmware/bin:$PATH"

#For MAC xterm errors  
export LC_ALL=en_US.UTF-8

